# wm/.local/share/systemd/user/corectrl.service
[Unit]
After=graphical-session.target

[Service]
Type=exec
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "wlroots:sway:Wayfire:labwc:Hyprland" ""
ExecStart=corectrl --minimize-systray
Restart=on-failure
Slice=app-graphical.slice

[Install]
WantedBy=graphical-session.target

